# templates/backups/b2-backup.service.j2
[Unit]
Description=Backup Unifi data to B2
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=simple
ExecStart=b2 sync {{ backup_directory }} {{ secrets.B2_DESTINATION_PATH }}
