# files/logdna/logdna.service
[Unit]
Description=LogDNA
RefuseManualStart=no
RefuseManualStop=no

[Service]
Type=simple
ExecStart=/usr/bin/logdna-agent
